# deploy/systemd/navi-offroute.service
[Unit]
Description=navi-offroute — off-network router + MVUM API (Echo6 navi-backend, extraction #8)
After=network-online.target
Wants=network-online.target

[Service]
# User=zvx is REQUIRED: PADUS PostGIS access uses libpq peer-auth via the local
# socket (dbname=padus, no password). zvx is verified to have padus access; a
# different user would fail to connect at runtime.
User=zvx
WorkingDirectory=/home/zvx/projects/repos/navi-backend
EnvironmentFile=/etc/navi-backend/navi-offroute.env
ExecStart=/home/zvx/projects/repos/navi-backend/.venv/bin/gunicorn 'services.navi_offroute.app:create_app()' --bind 127.0.0.1:8428 --workers 2 --timeout 130
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
