# deploy/systemd/navi-admin.service
[Unit]
Description=navi-admin — fleet admin-info aggregator (Echo6 navi-backend, extraction #7)
After=network-online.target
Wants=network-online.target

[Service]
User=zvx
WorkingDirectory=/home/zvx/projects/repos/navi-backend
EnvironmentFile=/etc/navi-backend/navi-admin.env
ExecStart=/home/zvx/projects/repos/navi-backend/.venv/bin/gunicorn 'services.navi_admin.app:create_app()' --bind 127.0.0.1:8427 --workers 2
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
