# deploy/systemd/navi-config.service
[Unit]
Description=navi-config — deployment profile API (Echo6 navi-backend, extraction #2)
After=network-online.target
Wants=network-online.target

[Service]
User=zvx
WorkingDirectory=/home/zvx/projects/repos/navi-backend
EnvironmentFile=/etc/navi-backend/navi-config.env
ExecStart=/home/zvx/projects/repos/navi-backend/.venv/bin/gunicorn 'services.navi_config.app:create_app()' --bind 127.0.0.1:8422 --workers 2
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
